# templates/logrotate.cron.erb
*/30  *  *  *  *   root    /usr/sbin/logrotate --verbose <%= @logrotate_vault_conf_file %>
